# conf.d/supervisor-ds.conf
[program:huaipeng_ds]
command=/srv/huaipeng/buildingdepot/venv/bin/uwsgi /etc/uwsgi/apps-available/huaipeng_ds.ini
  
directory=/srv/huaipeng/buildingdepot
autostart=true
autorestart=true
