# conf.d/supervisor-replica.conf
[program:rpc]
directory=/home/huaipeng/buildingdepot/CentralReplica
command=python main.py
