# configs/supervisor-oauth.conf
[program:oauth]
directory=/srv/buildingdepot/OAuth2Server
command=/srv/buildingdepot/venv/bin/python /srv/buildingdepot/OAuth2Server/app.py
