# configs/supervisor-cs.conf
[program:cs]
command=/srv/buildingdepot/venv/bin/uwsgi /etc/uwsgi/apps-available/cs.ini
  
directory=/srv/buildingdepot
autostart=true
autorestart=true
