# configs/supervisor-ds.conf
[program:ds]
command=/srv/buildingdepot/venv/bin/uwsgi /etc/uwsgi/apps-available/ds.ini
  
directory=/srv/buildingdepot
autostart=true
autorestart=true
