# configs/supervisor-replica.conf
[program:rpc]
directory=/srv/buildingdepot/CentralReplica
command=/srv/buildingdepot/venv/bin/python /srv/buildingdepot/CentralReplica/main.py
