# new_bd/conf.d/supervisor-cs.conf
[program:huaipeng_cs]
command=/srv/huaipeng/buildingdepot/venv/bin/uwsgi /etc/uwsgi/apps-available/huaipeng_cs.ini
  
directory=/srv/huaipeng/buildingdepot
autostart=true
autorestart=true
